# BBA_1.5_platform/targets/fs.dir/etc/inittab
::sysinit:/etc/init.d/rcS
::askfirst:/sbin/getty -L ttyS1 115200 vt100
